# server/deploy/pulse.service
[Unit]
Description=Pulse Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pulse
Group=pulse
ExecStart=/opt/pulse/pulse-server serve --config /opt/pulse/config.toml --data-dir /opt/pulse/data
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/pulse/data
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Allow binding to privileged ports
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
